# docker/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:keep-alive]
command= /home/run.sh
autorestart=true